# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Web server — command and user are injected via docker-compose environment
# (mirrors the official Laravel Sail supervisord.conf pattern)
[program:php]
command=%(ENV_SUPERVISOR_PHP_COMMAND)s
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Reverb WebSocket server — runs alongside the web server in the app container.
# Binding to 0.0.0.0 so it is reachable both from the host (browser)
# and from other containers (Horizon, Scheduler) on the Docker network.
[program:reverb]
command=php artisan reverb:start --host=0.0.0.0 --port=%(ENV_REVERB_PORT)s
directory=/var/www/html
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
